# setup/xinit.service
[Unit]
Description=xinit

[Service]
ExecStart=/bin/bash -c '! xset q && xinit'

[Install]
WantedBy=multi-user.target
